# brain-core/deploy/brain-core.service
[Unit]
Description=ONE-BRAIN core daemon
After=network-online.target nats-server.service
Wants=network-online.target
# NATS is not strictly required but preferred — daemon auto-falls back to HTTP-only.

[Service]
Type=simple
User=brain
Group=brain
WorkingDirectory=/srv/global-brain
Environment=NODE_ENV=production
Environment=BRAIN_DATA_DIR=/var/lib/brain
Environment=BRAIN_WORKSPACE=/srv/global-brain
Environment=BRAIN_NATS_URL=nats://127.0.0.1:4222
Environment=BRAIN_HTTP_HOST=0.0.0.0
Environment=BRAIN_HTTP_PORT=7070
Environment=BRAIN_CHECKPOINT_MS=60000
ExecStart=/usr/bin/node /srv/global-brain/brain-core/daemon.js
Restart=always
RestartSec=2
# Generous, brain caches can grow
MemoryHigh=2G
MemoryMax=4G
LimitNOFILE=65536
# Standard hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/brain /srv/global-brain

[Install]
WantedBy=multi-user.target
